# etc/s6-overlay/s6-rc.d/krb5kdc/run
#!/usr/bin/env bash
exec /usr/sbin/krb5kdc -n
