# etc/s6-overlay/s6-rc.d/kadmind/run
#!/usr/bin/env bash
exec /usr/sbin/kadmind -nofork
